# etc/com.isolator.docker-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.isolator.docker-link</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>target=$(readlink /var/run/docker.sock 2>/dev/null); if [ -n "$target" ] &amp;&amp; [ -S "$target" ]; then rm -f /var/run/docker.sock; ln "$target" /var/run/docker.sock; chown root:wheel /var/run/docker.sock; chmod 600 /var/run/docker.sock; fi</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>WatchPaths</key>
    <array>
        <string>/var/run/docker.sock</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/var/log/isolator-docker-link.log</string>
</dict>
</plist>
