# etc/com.isolator.docker-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.isolator.docker-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/socat</string>
        <string>UNIX-LISTEN:/var/run/docker-shared.sock,fork,group=staff,mode=770</string>
        <string>UNIX-CONNECT:/var/run/docker.sock</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/var/run/docker.sock</key>
            <true/>
        </dict>
    </dict>
    <key>StandardErrorPath</key>
    <string>/var/log/isolator-docker-proxy.log</string>
</dict>
</plist>
